# scripts/com.fatburn.autopush.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.fatburn.autopush</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/kejiawu/Library/Application Support/fatburn/auto-commit-push-logs.sh</string>
  </array>
  <!-- 07:55 local time (America/New_York on this Mac) — before 08:00 cloud email -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>7</integer>
    <key>Minute</key>
    <integer>55</integer>
  </dict>
  <key>RunAtLoad</key>
  <false/>
  <key>StandardOutPath</key>
  <string>/Users/kejiawu/Library/Logs/fatburn-autopush.launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/kejiawu/Library/Logs/fatburn-autopush.launchd.err.log</string>
</dict>
</plist>
